# b_infra/vm_oci-p-flex_1/1.os/idle-shutdown.timer
[Unit]
Description=Run idle shutdown check every 5 minutes
Documentation=https://github.com/diegonmarcos/back-System

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
AccuracySec=1min
Persistent=true

[Install]
WantedBy=timers.target
