# b_infra/vm_oci-p-flex_1/1.os/idle-shutdown.service
[Unit]
Description=Idle Shutdown Check
Documentation=https://github.com/diegonmarcos/back-System
After=docker.service

[Service]
Type=oneshot
ExecStart=/opt/scripts/idle-shutdown.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
